# quadlets/heimdall.container
[Unit]
Description=Heimdall Application Dashboard
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=heimdall
Image=lscr.io/linuxserver/heimdall:latest
Environment=TZ=America/New_York
Volume={{INSTALL_DIR}}/heimdall/config:/config:Z,U
Network=appnet
UserNS=auto

HealthCmd=["CMD-SHELL", "curl -f http://localhost:80/ || exit 1"]
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target