# quadlets/pwpush.container
[Unit]
Description=Password Pusher
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Container]
ContainerName=pwpush
Image=docker.io/pglombardo/pwpush:latest
Volume={{INSTALL_DIR}}/pwpush/storage:/opt/PasswordPusher/storage:Z,U
Network=appnet
# Note: Cannot use UserNS=auto - pwpush's Thruster proxy binds to port 80 internally

# RequiresPostgres:
Environment=DATABASE_URL=postgresql://{{PWPUSH_DB_USER}}:{{PWPUSH_DB_PASSWORD}}@postgresql:5432/{{PWPASS_DB_NAME}}
Environment=PWP__ENABLE_LOGINS=true
Environment=PWP__MAIL__RAISE_DELIVERY_ERRORS=false
Environment=PWP__HOST_DOMAIN={{PWPUSH_DOMAIN}}
Environment=PWP__HOST_PROTOCOL=https
Environment=PWP__ENABLE_FILE_PUSHES=true
Environment=PWP__FILES__STORAGE=local
Environment=PWP__FILES__EXPIRE_AFTER_DAYS_DEFAULT=2
Environment=PWP__FILES__EXPIRE_AFTER_DAYS_MAX=7
Environment=PWP__FILES__EXPIRE_AFTER_VIEWS_DEFAULT=5
Environment=PWP__FILES__EXPIRE_AFTER_VIEWS_MAX=10
Environment=PWP__FILES__RETRIEVAL_STEP_DEFAULT=true
Environment=PWP__ENABLE_URL_PUSHES=true
Environment=PWP__DISABLE_SIGNUPS=true

HealthCmd=["CMD-SHELL", "curl -f http://localhost:5100/ || exit 1"]
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target