# quadlets/rancher.container
[Unit]
Description=Rancher container
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=rancher
Image=docker.io/rancher/rancher:stable
Network=appnet
PodmanArgs=--privileged

# Persist Rancher data (includes embedded k3s/etcd)
Volume={{INSTALL_DIR}}/rancher:/var/lib/rancher:Z

# Rancher UI port (configure reverse proxy to point here)
PublishPort={{HOST_PORT}}:443

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
